# mirrmaid.cron
MAILTO=root
PATH=/sbin:/bin:/usr/sbin:/usr/bin
32 * * * * mirrmaid mirrmaid
